# pmacct/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pmacct/
drwxr-xr-x	root/root	etc/pmacct/examples/
-rw-r--r--	root/root	etc/pmacct/examples/agent_to_peer.map.example
-rw-r--r--	root/root	etc/pmacct/examples/allow-list.example
-rw-r--r--	root/root	etc/pmacct/examples/bgp_md5.lst.example
-rw-r--r--	root/root	etc/pmacct/examples/gnuplot-example.sh
-rw-r--r--	root/root	etc/pmacct/examples/gnuplot.script.example
-rw-r--r--	root/root	etc/pmacct/examples/is_symmetric.map.example
-rw-r--r--	root/root	etc/pmacct/examples/lpref.map.example
-rw-r--r--	root/root	etc/pmacct/examples/med.map.example
-rw-r--r--	root/root	etc/pmacct/examples/mrtg-example.sh
-rw-r--r--	root/root	etc/pmacct/examples/mrtg.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/networks.lst.example
-rw-r--r--	root/root	etc/pmacct/examples/nfacctd-print.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/nfacctd-sql_v1.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/nfacctd-sql_v2.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/peers.map.example
-rw-r--r--	root/root	etc/pmacct/examples/pmacctd-imt.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/pmacctd-multiple-plugins.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/pmacctd-sql_v1.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/pmacctd-sql_v2.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/pmacctd-sqlite3_v4.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/ports.lst.example
-rw-r--r--	root/root	etc/pmacct/examples/pretag.map.example
-rw-r--r--	root/root	etc/pmacct/examples/probe_netflow.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/probe_sflow.conf.example
-rw-r--r--	root/root	etc/pmacct/examples/rrdtool-example.sh
-rw-r--r--	root/root	etc/pmacct/nfacctd.conf.example
-rw-r--r--	root/root	etc/pmacct/pmacctd.conf.example
drwxr-xr-x	root/root	etc/pmacct/sql/
-rw-r--r--	root/root	etc/pmacct/sql/README.64bit
-rw-r--r--	root/root	etc/pmacct/sql/README.IPv6
-rw-r--r--	root/root	etc/pmacct/sql/README.agent_id2
-rw-r--r--	root/root	etc/pmacct/sql/README.cos
-rw-r--r--	root/root	etc/pmacct/sql/README.iface
-rw-r--r--	root/root	etc/pmacct/sql/README.is_symmetric
-rw-r--r--	root/root	etc/pmacct/sql/README.mask
-rw-r--r--	root/root	etc/pmacct/sql/README.mysql
-rw-r--r--	root/root	etc/pmacct/sql/README.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/README.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_bgp_v1.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v1.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v2.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v3.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v4.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v5.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v6.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v7.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-db_v8.mysql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_bgp_v1.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_bgp_v1.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v1.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v1.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v2.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v2.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v3.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v3.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v4.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v4.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v5.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v5.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v6.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v6.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v7.pgsql
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v7.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-create-table_v8.sqlite3
-rw-r--r--	root/root	etc/pmacct/sql/pmacct-grant-db.mysql
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/pmacct
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/pmacct
-rwxr-xr-x	root/root	usr/bin/pmmyplay
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/nfacctd
-rwxr-xr-x	root/root	usr/sbin/pmacctd
-rwxr-xr-x	root/root	usr/sbin/sfacctd
-rwxr-xr-x	root/root	usr/sbin/uacctd
