# quagga/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/quagga/
-rw-r--r--	root/root	etc/quagga/bgpd.conf.sample
-rw-r--r--	root/root	etc/quagga/bgpd.conf.sample2
-rw-r--r--	root/root	etc/quagga/ospfd.conf.sample
-rw-r--r--	root/root	etc/quagga/ripd.conf.sample
-rw-r--r--	root/root	etc/quagga/vtysh.conf.sample
-rw-r--r--	root/root	etc/quagga/zebra.conf.sample
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/quagga
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/vtysh
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/quagga/
-rw-r--r--	root/root	usr/include/quagga/buffer.h
-rw-r--r--	root/root	usr/include/quagga/checksum.h
-rw-r--r--	root/root	usr/include/quagga/command.h
-rw-r--r--	root/root	usr/include/quagga/distribute.h
-rw-r--r--	root/root	usr/include/quagga/filter.h
-rw-r--r--	root/root	usr/include/quagga/getopt.h
-rw-r--r--	root/root	usr/include/quagga/hash.h
-rw-r--r--	root/root	usr/include/quagga/if.h
-rw-r--r--	root/root	usr/include/quagga/if_rmap.h
-rw-r--r--	root/root	usr/include/quagga/jhash.h
-rw-r--r--	root/root	usr/include/quagga/keychain.h
-rw-r--r--	root/root	usr/include/quagga/linklist.h
-rw-r--r--	root/root	usr/include/quagga/log.h
-rw-r--r--	root/root	usr/include/quagga/md5.h
-rw-r--r--	root/root	usr/include/quagga/memory.h
-rw-r--r--	root/root	usr/include/quagga/memtypes.h
-rw-r--r--	root/root	usr/include/quagga/network.h
drwxr-xr-x	root/root	usr/include/quagga/ospfapi/
-rw-r--r--	root/root	usr/include/quagga/ospfapi/ospf_apiclient.h
drwxr-xr-x	root/root	usr/include/quagga/ospfd/
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_api.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_asbr.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_dump.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_ism.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_lsa.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_lsdb.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_nsm.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospf_opaque.h
-rw-r--r--	root/root	usr/include/quagga/ospfd/ospfd.h
-rw-r--r--	root/root	usr/include/quagga/plist.h
-rw-r--r--	root/root	usr/include/quagga/pqueue.h
-rw-r--r--	root/root	usr/include/quagga/prefix.h
-rw-r--r--	root/root	usr/include/quagga/privs.h
-rw-r--r--	root/root	usr/include/quagga/route_types.h
-rw-r--r--	root/root	usr/include/quagga/routemap.h
-rw-r--r--	root/root	usr/include/quagga/sigevent.h
-rw-r--r--	root/root	usr/include/quagga/smux.h
-rw-r--r--	root/root	usr/include/quagga/sockopt.h
-rw-r--r--	root/root	usr/include/quagga/sockunion.h
-rw-r--r--	root/root	usr/include/quagga/str.h
-rw-r--r--	root/root	usr/include/quagga/stream.h
-rw-r--r--	root/root	usr/include/quagga/table.h
-rw-r--r--	root/root	usr/include/quagga/thread.h
-rw-r--r--	root/root	usr/include/quagga/vector.h
-rw-r--r--	root/root	usr/include/quagga/version.h
-rw-r--r--	root/root	usr/include/quagga/vty.h
-rw-r--r--	root/root	usr/include/quagga/workqueue.h
-rw-r--r--	root/root	usr/include/quagga/zassert.h
-rw-r--r--	root/root	usr/include/quagga/zclient.h
-rw-r--r--	root/root	usr/include/quagga/zebra.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libospf.a
-rwxr-xr-x	root/root	usr/lib/libospf.la
lrwxrwxrwx	root/root	usr/lib/libospf.so -> libospf.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libospf.so.0 -> libospf.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libospf.so.0.0.0
-rw-r--r--	root/root	usr/lib/libospfapiclient.a
-rwxr-xr-x	root/root	usr/lib/libospfapiclient.la
lrwxrwxrwx	root/root	usr/lib/libospfapiclient.so -> libospfapiclient.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libospfapiclient.so.0 -> libospfapiclient.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libospfapiclient.so.0.0.0
-rw-r--r--	root/root	usr/lib/libzebra.a
-rwxr-xr-x	root/root	usr/lib/libzebra.la
lrwxrwxrwx	root/root	usr/lib/libzebra.so -> libzebra.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libzebra.so.0 -> libzebra.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libzebra.so.0.0.0
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/bgpd
-rwxr-xr-x	root/root	usr/sbin/ospfclient
-rwxr-xr-x	root/root	usr/sbin/ospfd
-rwxr-xr-x	root/root	usr/sbin/ripd
-rwxr-xr-x	root/root	usr/sbin/watchquagga
-rwxr-xr-x	root/root	usr/sbin/zebra
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	quagga/quagga	var/run/quagga/
